# instance-b/config/supervisord.conf
[supervisord]
user=root
logfile = /tmp/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /tmp/supervisord.pid

[program:sshd]
command=/usr/sbin/sshd -D

[program:python]
command=/usr/bin/python3 /app/src/main.py
autostart=true
autorestart=true
startsecs=5
